# examples/gitwatch@.service
[Unit]
Description=Watch file or directory and git commit all changes. run with: systemctl --user --now enable gitwatch@$(systemd-escape -- "'-r url/to/repository' /path/to/folder").service
# Add dependencies for networking, as -r implies network access
After=network-online.target
Wants=network-online.target

[Service]
Environment="SCRIPT_ARGS=%I"
# Add coreutils and util-linux to the PATH
# This ensures git, flock, and timeout are found.
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
# Specify dependency packages for clarity (though PATH is the fix)
# Requires=git, util-linux (flock), coreutils (timeout), inotify-tools, procps (pkill)
ExecStart=/usr/bin/bash -c "/usr/local/bin/gitwatch ${SCRIPT_ARGS}"
ExecStop=/bin/true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
